# tvOS/Package/Library/LaunchDaemons/com.matchstic.reprovisiond.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.matchstic.reprovisiond</string>
	<key>MachServices</key>
	<dict>
		<key>com.matchstic.reprovisiond</key>
		<true/>
	</dict>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/bin/reprovisiond</string>
	</array>
	<key>UserName</key>
	<string>mobile</string>
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<true/>
	<key>StandardOutPath</key>
	<string>/var/mobile/Library/Logs/reprovisiond.log</string>
	<key>StandardErrorPath</key>
	<string>/var/mobile/Library/Logs/reprovisiond-Error.log</string>
</dict>
</plist>
